# crontab
# Update libc-database every day
# min hour day mon week  user  cmd
  0   12   *   *   *     root  cd /libc-database && (./get; echo -n "### Update done "; date) >> /var/log/libcdb.log 2>&1
