# cron.conf
[program:cron]
command = /usr/sbin/cron -f
user = root
autostart = true

[program:cron_output]
command = tail -f /var/log/libcdb.log
autostart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
